# deploy/systemd/rssapp-staging.timer
[Unit]
Description=Check GHCR for a newer rssapp staging image

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
